# profile.d/experimental
JENKINS_GIT_BRANCH=master
JENKINS_GIT_REPOSITORY=git@github.com:olblak/jenkins.git
GIT_EMAIL=release@jenkins.io
GIT_NAME=release-bot
GPG_FILE=gpg-test-jenkins-release.gpg
GPG_KEYNAME=test-jenkins-release
MAVEN_REPOSITORY_URL='http://nexus/repository'
MAVEN_REPOSITORY_NAME=maven-releases
MAVEN_PUBLIC_JENKINS_REPOSITORY_MIRROR_URL='http://nexus/repository/jenkins-public/'
SIGN_ALIAS=jenkins
