# profile.d/experimental
RELEASE_GIT_BRANCH=master
RELEASE_GIT_REPOSITORY=git@github.com:olblak/jenkins.git
JENKINS_VERSION=latest
GIT_EMAIL=66998184+jenkins-release-bot@users.noreply.github.com
GIT_NAME="Jenkins Release Bot"
GPG_KEYNAME="5E386EADB55F01504CAE8BCF7198F4B714ABFC68"
GPG_VAULT_NAME="jenkins-release-pgp-2026"
MAVEN_REPOSITORY_URL='https://repo.jenkins-ci.org'
MAVEN_REPOSITORY_NAME=olblak-sandbox
MAVEN_PUBLIC_JENKINS_REPOSITORY_MIRROR_URL='https://repo.jenkins-ci.org/public/'
SIGN_ALIAS=jenkins

# PACKAGING ENV used from jenkinsci/packaging
RELEASELINE="-experimental"
