# profile.d/experimental
RELEASE_GIT_BRANCH=master
RELEASE_GIT_REPOSITORY=git@github.com:olblak/jenkins.git
JENKINS_VERSION=latest
GIT_EMAIL=66998184+jenkins-release-bot@users.noreply.github.com
GIT_NAME="Jenkins Release Bot"
GPG_KEYNAME="37308192E79DCE3E3833A32194FD886DA98761AC"
GPG_VAULT_NAME="jenkins-release-pgp-2023"
MAVEN_REPOSITORY_URL='https://repo.jenkins-ci.org'
MAVEN_REPOSITORY_NAME=olblak-sandbox
MAVEN_PUBLIC_JENKINS_REPOSITORY_MIRROR_URL='https://repo.jenkins-ci.org/public/'
SIGN_ALIAS=jenkins

# PACKAGING ENV used from jenkinsci/packaging
RELEASELINE="-experimental"

BUILDENV=env/azure.mk

# Define endpoint used to force mirror synchronization
PKGSERVER=mirrorbrain@20.186.155.37
PKGSERVER_SSH_OPTS="-p 22 -o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null"
