# profile.d/experimental
JENKINS_GIT_BRANCH=master
JENKINS_GIT_REPOSITORY=git@github.com:olblak/jenkins.git
GIT_EMAIL=release@jenkins.io
GIT_NAME=release-bot
GPG_FILE=jenkins-release.gpg
GPG_KEYNAME=jenkins-release
MAVEN_REPOSITORY_URL='https://repo.jenkins-ci.org'
MAVEN_REPOSITORY_NAME=olblak-sandbox
MAVEN_PUBLIC_JENKINS_REPOSITORY_MIRROR_URL='http://nexus/repository/jenkins-public/'
SIGN_ALIAS=jenkins
