# profile.d/security
#
# WARNING: Any variables defined here, override those defined from the Jenkinsfile
#
# RELEASE_GIT_BRANCH is provided at run time
#
RELEASE_GIT_REPOSITORY=git@github.com:jenkinsci-cert/jenkins.git
#
# JENKINS_VERSION is provided at run time
#
GIT_EMAIL=66998184+jenkins-release-bot@users.noreply.github.com
GIT_NAME="Jenkins Release Bot"
GPG_KEYNAME="63667EE74BBA1F0A08A698725BA31D57EF5975CA"
GPG_VAULT_NAME="jenkins-release-pgp-2023"
MAVEN_REPOSITORY_URL=https://repo.jenkins-ci.org
MAVEN_PUBLIC_JENKINS_REPOSITORY_MIRROR_URL='https://repo.jenkins-ci.org/public/'
SIGN_ALIAS=jenkins

### Git Staging Promotion Settings
#
RELEASE_GIT_PRODUCTION_REPOSITORY=git@github.com:jenkinsci/jenkins.git
RELEASE_GIT_PRODUCTION_BRANCH=master

### Maven Staging Promotion Settings
#
# Remark: Version to be promoted will be the latest version pushed to the maven repository
#         as retrieved from ./utils/getJenkinsVersion.py

MAVEN_REPOSITORY_PRODUCTION_NAME=releases

RELEASE_GIT_BRANCH=security-master
MAVEN_REPOSITORY_NAME=chameleone
JENKINS_VERSION=2.442
RELEASELINE=
