# profile.d/security
#
# WARNING: Any variables defined here, override those defined from the Jenkinsfile
#
# RELEASE_GIT_BRANCH is provided at run time
#
RELEASE_GIT_REPOSITORY=git@github.com:jenkinsci-cert/jenkins.git
#
# JENKINS_VERSION is provided at run time
#
GIT_EMAIL=66998184+jenkins-release-bot@users.noreply.github.com
GIT_NAME="Jenkins Release Bot"
GPG_KEYNAME=62A9756BFD780C377CF24BA8FCEF32E745F2C3D5
GPG_VAULT_NAME=jenkins-release-pgp
MAVEN_REPOSITORY_URL=https://repo.jenkins-ci.org
MAVEN_PUBLIC_JENKINS_REPOSITORY_MIRROR_URL='https://repo.jenkins-ci.org/public/'
SIGN_ALIAS=jenkins

### Git Staging Promotion Settings
#
RELEASE_GIT_PRODUCTION_REPOSITORY=git@github.com:jenkinsci/jenkins.git
RELEASE_GIT_PRODUCTION_BRANCH=master

### Maven Staging Promotion Settings
#
# Remark: Version to be promoted will be the latest version pushed to the maven repository
#         as retrieved from ./utils/getJenkinsVersion.py

MAVEN_REPOSITORY_PRODUCTION_NAME=releases

### EVERYTHING ABOVE IS UNUSED -- NO PUBLISHING OF SECURITY FIXES INSIDE THE RELEASE FLOW
MAVEN_REPOSITORY_NAME=camelia
RELEASE_GIT_BRANCH=security-stable-2.319
JENKINS_VERSION=2.319.3
RELEASELINE=-stable
